# crontab-task/supervisord.conf
[program:flaskdemo]
directory=D:\task
command=D:\\BinConfig\\python310\\python.exe vad.py
stdout_logfile=D:\\task\\flask.log
numprocs=1
autostart=true
autorestart=true
redirect_stderr = true
